# supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
command=gunicorn /deploy/myblog/myblog/wsgi:application
autorestart=true
